# cron/cron.d/backup
SHELL=/bin/sh

14 2 * * * root lftp -f /backup/backup.lftp
