# cron/cron.d/backup
SHELL=/bin/sh

42 17 * * * root lftp -f /backup/backup.lftp
